# launchd/com.automation.telegram-relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.automation.telegram-relay</string>
	<key>ProgramArguments</key>
	<array>
		<string>__PYTHON3__</string>
		<string>__AUTOMATION_DIR__/telegram-relay.py</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>WorkingDirectory</key>
	<string>__AUTOMATION_DIR__</string>
	<key>StandardOutPath</key>
	<string>__HOME__/logs/relay/launchd_stdout.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/logs/relay/launchd_stderr.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>__HOME__</string>
	</dict>
	<key>ThrottleInterval</key>
	<integer>5</integer>
</dict>
</plist>
