# launchd/com.automation.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.automation.heartbeat</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>__AUTOMATION_DIR__/heartbeat.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>1800</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>WorkingDirectory</key>
	<string>__AUTOMATION_DIR__</string>
	<key>StandardOutPath</key>
	<string>__HOME__/logs/heartbeat/launchd_stdout.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/logs/heartbeat/launchd_stderr.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>__HOME__</string>
	</dict>
</dict>
</plist>
